# delete_old_rec.service
[Unit]
Description=Delete old camera recordings

[Service]
ExecStart=/usr/bin/find /srv/http/rec/event/ -mtime +30 -delete

